# apolo_server/conf/supervisor/backend.ini
[program:backend]
process_name=%(program_name)s_%(process_num)02d
command=uwsgi --reload /tmp/uwsgi/uwsgi.pid
directory=/home/apolo/project
autostart=false
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/home/apolo/project/apolo/apolo_server/logs/backend.log