# apolo_server/conf/supervisor/frontend.ini
[program:frontend]
process_name=%(program_name)s_%(process_num)02d
command=npm restart
directory=/home/apolo/project/apolo/frontend/APP
autostart=false
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/home/apolo/project/apolo/apolo_server/logs/frontend.log