# apolo_server/conf/supervisor/parser_worker.ini
[program:parser_worker]
process_name=%(program_name)s_%(process_num)02d
command=nohup /home/apolo/env/bin/python /home/apolo/project/apolo/apolo_server/processor/worker/parser_worker.py
directory=/home/apolo/project/apolo/apolo_server/processor/worker
autostart=false
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/home/apolo/project/apolo/apolo_server/logs/parser_process.log