# deploy/supervisor/backend.ini
[program:backend]
process_name=%(program_name)s_%(process_num)02d
command=/home/apolo/env/bin/uwsgi --ini /home/apolo/project/apolo/deploy/uwsgi/uwsgi.ini
directory=/home/apolo/project/apolo
autostart=false
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/home/apolo/project/apolo/logs/backend_process.log