# deploy/supervisor/cli_worker.ini
[program:cli_worker]
process_name=%(program_name)s_%(process_num)02d
command=nohup /home/apolo/env/bin/python /home/apolo/project/apolo/apolo_server/processor/worker/cisco_cli_worker.py
directory=/home/apolo/project/apolo/apolo_server/processor/worker
autostart=false
user=apolo
numprocs=5
redirect_stderr=true
stdout_logfile=/home/apolo/project/apolo/logs/cli_process.log