# deploy/supervisor/status_snmp_worker.ini
[program:status_snmp_worker]
process_name=%(program_name)s_%(process_num)02d
command=nohup /home/apolo/env/bin/python /home/apolo/project/apolo/apolo_server/processor/worker/snmp_status_worker.py
directory=/home/apolo/project/apolo/apolo_server/processor/worker
autostart=false
user=apolo
numprocs=10
redirect_stderr=true
stdout_logfile=/home/apolo/project/apolo/logs/snmp_status_process.log